# systemd/complex-frogs.service
[Unit]
Description=FastAPI web server for the Complex Frogs project
# start after the network is up
After=network.target

[Service]
ExecStart=/path/to/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
WorkingDirectory=/path/to/app
User=username
Group=group
Restart=always

[Install]
WantedBy=multi-user.target